# srcpkgs/hal/files/hal.rc
#!/sbin/runscript

command=/usr/sbin/hald
pidfile=/var/run/hald/pid
name="HAL (Hardware Abstraction Layer)"

depend()
{
	need dbus
}
