# srcpkgs/bluez/files/bluetooth.service
[Unit]
Description=Bluetooth Manager

[Service]
Type=dbus
BusName=org.bluez
ExecStart=/usr/sbin/bluetoothd -n

[Install]
WantedBy=bluetooth.target
