# srcpkgs/systemd-units/files/acpid.service
[Unit]
Description=ACPI Event Daemon
Requires=acpid.socket

[Service]
EnvironmentFile=/etc/conf.d/acpid
ExecStart=/usr/sbin/acpid -f ${ACPID_ARGS}
