# srcpkgs/systemd-units/files/nscd.socket
[Unit]
Description=Name Service Cache Daemon Socket

[Socket]
ListenStream=/var/run/nscd/socket

[Install]
WantedBy=sockets.target
