# srcpkgs/ufw/files/ufw.service
[Unit]
Description=CLI Netfilter Manager

[Service]
Type=oneshot
ExecStart=/lib/ufw/ufw-init start
ExecStop=/lib/ufw/ufw-init stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
