# srcpkgs/ufw/files/ufw.rc
#!/sbin/runscript

depend() {
	before net
}

start() {
	ebegin "Starting ufw"
	ufw enable
	eend $?
}

stop() {
	ebegin "Stopping ufw"
	ufw disable
	eend $?
}
