# srcpkgs/ntp/files/ntpd.service
[Unit]
Description=Network Time Service

[Service]
ExecStart=/usr/bin/ntpd -ng

[Install]
WantedBy=multi-user.target
