# srcpkgs/opensmtpd/files/smtpd.service
[Unit]
Description=OpenSMTPD
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/smtpd

[Install]
WantedBy=multi-user.target
