# srcpkgs/opensmtpd/files/smtpd.socket
[Unit]
Description=OpenSMTPD Socket

[Socket]
ListenStream=/run/smtpd.sock

[Install]
WantedBy=sockets.target
