# templates/nfs-utils/files/nfslock.rc
#!/sbin/runscript

command=/usr/sbin/rpc.statd
command_args=$nfslock_args
pidfile=/var/run/rpc.statd.pid
name="RPC Statd daemon"

depend()
{
	need portmap
}
