# templates/dbus/files/dbus
#!/bin/sh
#
# PROVIDE: dbus
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name="dbus"
rcvar=$name
command="/usr/bin/dbus-daemon"
command_args="--system"
pidfile="/var/run/dbus/pid"
start_precmd=dbus_prestart

dbus_prestart() {
	dir="/var/run/dbus"
	if [ ! -d $dir ]; then
		mkdir $dir
		chmod 0755 $dir
		chown dbus:dbus $dir
	fi
	/usr/bin/dbus-uuidgen --ensure
}

load_rc_config $name
run_rc_command "$1"
