# templates/hal/files/hal
#!/bin/sh
#
# PROVIDE: hal
# REQUIRE: DAEMON dbus
# KEYWORD: shutdown

. /etc/rc.subr

name="hal"
rcvar=$name
command="/usr/sbin/hald"
pidfile="/var/run/hald/pid"
start_precmd=hal_prestart

hal_prestart() {
	dirs="/var/run/hald /var/cache/hald"
	for d in ${dirs}; do
		if [ ! -d ${d} ]; then
			mkdir ${d}
			chmod 0755 ${d}
			chown hal:hal ${d}
		fi
	done
	for f in hald-local hald-runner; do
		if [ ! -x /var/run/hald/${f} ]; then
			touch /var/run/hald/${f}
		fi
	done
}

load_rc_config $name
run_rc_command "$1"
