# srcpkgs/cyrus-sasl/files/saslauthd.rc
#!/sbin/runscript

command=/usr/sbin/saslauthd
command_args="${SASLAUTHD_ARGS}"

depend()
{
	need net
}
