# srcpkgs/dcron/files/crond.init
#!/sbin/runscript

command=/usr/sbin/crond
command_args="${CROND_ARGS}"

depend()
{
	need localmount
	after bootmisc
}
