# srcpkgs/dcron/files/dcron.service
[Unit]
Description=Periodic Command Scheduler

[Service]
Type=forking
ExecStart=/usr/sbin/crond
Restart=always

[Install]
WantedBy=multi-user.target
