# srcpkgs/dcron/files/dcron.service
[Unit]
Description=Periodic Command Scheduler

[Service]
ExecStart=/usr/sbin/crond -f

[Install]
WantedBy=multi-user.target
