# srcpkgs/samba/files/nmbd.service
[Unit]
Description=Samba NetBIOS name server

[Service]
ExecStart=/usr/sbin/nmbd -F

[Install]
WantedBy=multi-user.target
