# srcpkgs/samba/files/smbd@.service
[Unit]
Description=Samba SMB/CIFS server instance

[Service]
ExecStart=/usr/sbin/smbd -F
ExecReload=/usr/bin/kill -HUP $MAINPID
StandardInput=socket
