# srcpkgs/samba/files/smbd.socket
[Unit]
Description=Samba SMB/CIFS server socket

[Socket]
ListenStream=445
Accept=yes

[Install]
WantedBy=sockets.target
