# srcpkgs/cntlm/files/cntlm.service
[Unit]
Description=Cntlm Proxy Server
After=network.target

[Service]
ExecStart=/usr/bin/cntlm -f

[Install]
WantedBy=multi-user.target
