# srcpkgs/lightdm/files/lightdm.rc
#!/sbin/runscript

command=/usr/bin/lightdm
name="Light Display Manager"

depend()
{
	need dbus localmount
}
