# srcpkgs/x13s-base/files/72-touchscreen.rules
ACTION=="add" , RUN+="/bin/sh -c 'echo 4-0010 > /sys/bus/i2c/drivers/i2c_hid_of/bind'"
