# srcpkgs/x13s-base/files/90-net-address.rules
ACTION=="add", SUBSYSTEM=="net", KERNELS=="0006:01:00.0", RUN+="/usr/libexec/x13s-setup wlan $name"
