# srcpkgs/rpcbind/files/rpcbind.service
[Unit]
Description=RPC Bind

[Service]
Type=forking
ExecStart=/usr/bin/rpcbind

[Install]
WantedBy=multi-user.target
