# srcpkgs/sane/files/saned.socket
[Unit]
Description=saned incoming socket

[Socket]
ListenStream=6566
Accept=yes

[Install]
WantedBy=sockets.target
