# srcpkgs/sane/files/saned@.service
[Unit]
Description=Scanner Service

[Service]
Group=scanner
ExecStart=/usr/sbin/saned
StandardInput=socket
StandardError=syslog
