# srcpkgs/connman/files/connmand.rc
#!/sbin/runscript

command=/usr/sbin/connmand
name="Connection Manager daemon"

depend()
{
	need dbus localmount
}
