# srcpkgs/dnsmasq/files/dnsmasq.rc
#!/sbin/runscript

pidfile=/var/run/dnsmasq.pid
command=/usr/sbin/dnsmasq
command_args="${DNSMASQ_OPTS:--g nogroup -u nobody}"

depend()
{
	provide dns
	need localmount net
	after bootmisc
}
