# srcpkgs/anacron/files/anacrond.rc
#!/sbin/runscript
#

depend()
{
	use logger
	need clock hostname 
}

start()
{
	ebegin "Starting anacron"
	/usr/sbin/anacron -s >> /var/log/cron 2>&1
	eend $?
}
