# srcpkgs/chrony/files/chrony.service
[Unit]
Description=Chrony Network Time Daemon

[Service]
ExecStart=/usr/sbin/chronyd -n -u chrony

[Install]
WantedBy=multi-user.target
