# srcpkgs/acpid/files/acpid.socket
[Unit]
Description=ACPID Listen Socket

[Socket]
ListenStream=/run/acpid.socket

[Install]
WantedBy=sockets.target
