# srcpkgs/portmap/files/portmap.rc
#!/sbin/runscript

: ${portmap_args:=-t /var/chroot/portmap}

command=/sbin/portmap
command_args=$portmap_args
name="RPC Port mapper daemon"

depend()
{
	need net
}
