# srcpkgs/odroid-u2-base/files/odroid-led.service
[Unit]
Description=Sets Odroids LED on boot

[Service]
Type=oneshot
ExecStart=/usr/sbin/odroid-led

[Install]
WantedBy=multi-user.target
