# srcpkgs/inetutils/files/ftpd.service
[Unit]
Description=FTPD Daemon
After=network.target

[Service]
ExecStart=/usr/sbin/ftpd -D
Type=forking

[Install]
WantedBy=multi-user.target
