# srcpkgs/tinc/files/tincd@.service
[Unit]
Description=Tinc daemon for net %i
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/tincd -n %i

[Install]
WantedBy=multi-user.target
