# srcpkgs/squid/files/cron.daily
#!/bin/sh
if sv status squid | grep -q '^run'; then
	squid -k rotate
fi
