# srcpkgs/mdocml/files/makewhatis.cron.daily
#!/bin/sh
/sbin/makewhatis -a -Tutf8
exit 0
