# srcpkgs/wicd/files/wicd.rc
#!/sbin/runscript

command=/usr/sbin/wicd
pidfile=/var/run/wicd/wicd.pid

depend()
{
	need dbus
	need hal
}
