# srcpkgs/nfs-utils/files/proc-fs-nfsd.mount
[Unit]
Description=NFS Daemon filesystem
DefaultDependencies=no

[Mount]
What=sunrpc
Where=/proc/fs/nfsd
Type=nfsd
