# srcpkgs/nfs-utils/files/proc-fs-nfsd.mount
[Unit]
Description=RPC Pipe File System
DefaultDependencies=no

[Mount]
What=sunrpc
Where=/proc/fs/nfsd
Type=rpc_pipefs
