# srcpkgs/nfs-utils/files/rpc-statd.service
[Unit]
Description=NFSv2/3 Network Status Monitor Daemon
After=rpcbind.service
Requires=rpcbind.service

[Service]
Type=forking
ExecStart=/usr/sbin/rpc.statd

[Install]
WantedBy=multi-user.target
