# srcpkgs/nfs-utils/files/var-lib-nfs-rpc_pipefs.mount
[Unit]
Description=RPC pipe filesystem

[Mount]
What=rpc_pipefs
Where=/var/lib/nfs/rpc_pipefs
Type=rpc_pipefs
