# srcpkgs/dhcpcd/files/dhcpcd@.service
[Unit]
Description=Lightweight DHCP client daemon on %I
Wants=network.target
Before=network.target

[Service]
ExecStart=/usr/sbin/dhcpcd -q -B %I
ExecStop=/usr/sbin/dhcpcd -k %I

[Install]
Alias=multi-user.target.wants/dhcpcd@eth0.service
