# srcpkgs/fake-hwclock/files/rc.local
if (command -v /sbin/fake-hwclock >/dev/null 2>&1) ; then
	[ -r /etc/default/fake-hwclock ] && . /etc/default/fake-hwclock
	/sbin/fake-hwclock load $FORCE
fi
