# srcpkgs/cronie/files/crond.rc
#!/sbin/runscript

command=/usr/sbin/crond
pidfile=/var/run/crond.pid

depend()
{
	need localmount
	after bootmisc
}
