# templates/initscripts/files/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#


