# templates/initscripts/files/rc.multi
#!/bin/sh
#
# This script is called when runlevel 3 (multiuser) is entered.
#

export HOME=/
export PATH=/sbin:/bin:/usr/sbin:/usr/bin
umask 022

. /etc/rc.subr
. /etc/rc.conf
_rc_conf_loaded=true

# Start services
files=$(rcorder -s nostart ${rc_rcorder_flags} /etc/rc.d/*)

for _rc_elem in $files; do
	run_rc_script $_rc_elem start
done

[ -x /etc/rc.local ] && /etc/rc.local
